# _lib/service/list/run
#!/usr/bin/env bash

kproject chart list

