# commands/service/list/run
#!/usr/bin/env bash

set -eo pipefail

kproject chart list

