# commands/service/list/run
#!/usr/bin/env bash

set -eo pipefail

(
  command cd "$PROJECT_ROOT_PATH/services"
  sed -n '/include/,${s/.*name: *//p}' < "${ENVIRONMENT_NAME}-cluster.yaml"
)

