# commands/service/deploy/run
#!/usr/bin/env bash

set -eo pipefail

INCLUDE=$(echo "$@" | $SED 's/ /,/g')
(
  command cd "$PROJECT_ROOT_PATH/services"
  kontemplate apply "${ENVIRONMENT_NAME}-cluster.yaml" --include=$INCLUDE
)

